# bubba-backend/debian/bubba-backend.bubba-notify.cron.d
*/5 * * * * root test -x /usr/lib/web-admin/notify-dispatcher.pl && /usr/lib/web-admin/notify-dispatcher.pl
